# studentbuddy.service
[Unit]
Description=StudentBuddy Telegram Bot
After=network.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/studentbuddy
Environment="PATH=/opt/studentbuddy/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Загрузка переменных окружения из .env
EnvironmentFile=/opt/studentbuddy/.env

# Запуск бота
ExecStart=/opt/studentbuddy/venv/bin/python3 /opt/studentbuddy/bot.py

# Автоматический рестарт при падении
Restart=always
RestartSec=10

# Ограничения ресурсов
MemoryLimit=512M
CPUQuota=100%

# Логирование
StandardOutput=append:/opt/studentbuddy/logs/bot.log
StandardError=append:/opt/studentbuddy/logs/bot_error.log

# Безопасность
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/studentbuddy/data /opt/studentbuddy/logs

[Install]
WantedBy=multi-user.target
